# scripts/launchagents/ai.krab.weekly-heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.krab.weekly-heartbeat</string>

    <key>Comment</key>
    <string>Wave 114: еженедельный owner heartbeat DM с сводкой. Понедельник 09:30 local.</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin/python</string>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/scripts/krab_weekly_heartbeat.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/weekly_heartbeat.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/weekly_heartbeat.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/pablito</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>KRAB_PANEL_URL</key>
        <string>http://127.0.0.1:8080</string>
    </dict>
</dict>
</plist>
