# scripts/launchagents/ai.krab.google-sales-followup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.krab.google-sales-followup</string>

    <key>Comment</key>
    <string>Wave 117: ежедневный reminder про Google Sales POC по Anthropic Vertex quota. Daily 10:00 local — после morning brief и Wave 104 preflight (09:00).</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin/python</string>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/scripts/krab_google_sales_followup.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/google_sales_followup.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/google_sales_followup.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/pablito</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
