# scripts/launchagents/ai.krab.dependency-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Wave 110: Dependency Vulnerability Audit (weekly, FREE)

        Запускает scripts/krab_dependency_audit.py каждую субботу в 04:00
        (низкая нагрузка). Результат:
          - JSON report в ~/.openclaw/krab_runtime_state/dep_audit.json (last 10 runs)
          - Prometheus Gauge krab_dependency_vulns_total{severity}
          - Alert CriticalDependencyVuln при critical > 0 за 1h

        Setup:
          venv/bin/pip install pip-audit

        Install:
          launchctl bootstrap gui/501 \
            /Users/pablito/Antigravity_AGENTS/Краб/scripts/launchagents/ai.krab.dependency-audit.plist
    -->
    <key>Label</key>
    <string>ai.krab.dependency-audit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin/python</string>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/scripts/krab_dependency_audit.py</string>
    </array>

    <!-- Weekly: Saturday 04:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/pablito</string>
        <key>PATH</key>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/pablito/Antigravity_AGENTS/Краб</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/dep_audit.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/dep_audit.err.log</string>

    <key>KeepAlive</key>
    <false/>

    <key>LowPriorityIO</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
