# scripts/launchagents/ai.krab.backup-verify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Wave 105: Backup integrity verification (LaunchAgent)

        Weekly SHA256 + sqlite3 PRAGMA integrity_check для всех backup
        файлов (sessions, runtime state, memory archive). Каждое
        воскресенье 05:00 — после nightly maintenance, перед secrets
        audit (09:00).

        Install:
          launchctl bootstrap gui/501 \
            /Users/pablito/Antigravity_AGENTS/Краб/scripts/launchagents/ai.krab.backup-verify.plist
    -->
    <key>Label</key>
    <string>ai.krab.backup-verify</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin/python</string>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/scripts/krab_backup_verify.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/pablito</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/backup_verify.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/backup_verify.err.log</string>

    <key>KeepAlive</key>
    <false/>

    <key>LowPriorityIO</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
