# scripts/launchagents/ai.krab.inbox-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.krab.inbox-watcher</string>

    <key>Comment</key>
    <string>Krab Inbox Watcher — monitors ~/Krab_Inbox for new files</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin/python</string>
      <string>/Users/pablito/Antigravity_AGENTS/Краб/scripts/krab_inbox_watcher.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб</string>

    <key>StandardOutPath</key>
    <string>/tmp/krab-inbox-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/krab-inbox-watcher.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/pablito</string>
      <key>PATH</key>
      <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
