# scripts/launchagents/ai.krab.core.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.krab.core</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin/python3</string>
    <string>-m</string>
    <string>src.main</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/pablito/Antigravity_AGENTS/Краб</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>OPENCLAW_BROWSER_AUTOSTART</key>
    <string>1</string>
    <key>DEDICATED_CHROME_ENABLED</key>
    <string>true</string>
    <key>DEDICATED_CHROME_PORT</key>
    <string>9222</string>
    <key>DEDICATED_CHROME_PROFILE_DIR</key>
    <string>/tmp/krab-chrome</string>
    <!-- Wave 39-B (2026-05-09): proactive event detection (joins/media/ai). -->
    <key>KRAB_PROACTIVE_ENABLED</key>
    <string>1</string>
    <!-- Wave 44-M (2026-05-09): bump wall-clock caps 600s→1800s + 300s→1800s
         для длинных agentic tasks (research, multi-tool chains, codex-cli).
         Liveness check (Wave 36-A zombie + Wave 44-I active probe) — отдельный
         safety net через update_id frozen detection. -->
    <key>KRAB_LLM_WALL_CLOCK_CAP_SEC</key>
    <string>1800</string>
    <key>KRAB_LLM_WALL_CLOCK_CAP_PHOTO_SEC</key>
    <string>1800</string>
    <key>KRAB_RESPONSE_HARD_CAP_SEC</key>
    <string>1800</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>Crashed</key>
    <true/>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <!-- ThrottleInterval=60: разрывает respawn-storm при corruption loop -->
  <!-- (incident 26.04.2026 — 338 рестартов / 5.5h). Минимум 1 restart/min. -->
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/krab_launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/pablito/Antigravity_AGENTS/Краб/logs/krab_launchd.err.log</string>
</dict>
</plist>
