# scripts/launchagents/ai.krab.inbox-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.krab.inbox-cleanup</string>

    <key>Comment</key>
    <string>Wave 34-C + Wave 85: Ежедневная авточистка inbox в два этапа. Этап 1 (Wave 85): bulk-ack stale open proactive_action items старше INBOX_CLEANUP_BULK_ACK_HOURS (12h). Этап 2 (34-C): архивация open items старше INBOX_CLEANUP_MAX_AGE_DAYS (7d) из безопасных kinds. Запускается каждый день в 04:00.</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin/python</string>
        <string>/Users/pablito/Antigravity_AGENTS/Краб/scripts/inbox_cleanup_stale.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/krab_inbox_cleanup.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/krab_inbox_cleanup.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/pablito</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>KRAB_PANEL_URL</key>
        <string>http://127.0.0.1:8080</string>
        <key>INBOX_CLEANUP_MAX_AGE_DAYS</key>
        <string>7</string>
        <key>INBOX_CLEANUP_BULK_ACK_HOURS</key>
        <string>12</string>
        <key>INBOX_CLEANUP_BULK_ACK_KIND</key>
        <string>proactive_action</string>
        <key>INBOX_CLEANUP_BULK_ACK_ENABLED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
