# scripts/launchagents/com.krab.mcp-tor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  DEPRECATED Wave 50-B (2026-05-10): replaced by tor-full (25 tools).
  Plist preserved для git history. LaunchAgent больше НЕ должен загружаться:
  оркестратор выполняет `launchctl unload ~/Library/LaunchAgents/com.krab.mcp-tor.plist`
  и `rm` копии из ~/Library/LaunchAgents/. Replacement — `tor-full` MCP
  (`/Users/pablito/Antigravity_AGENTS/tor-mcp/server.py`, stdio, 25 tools)
  обладает superset функциональностью (browse_onion, circuit/identity мгмт,
  shodan/censys/virustotal, gpg, port_scan и др).

  Этот файл оставлен в репо как:
  - reference SSE LaunchAgent template для будущих MCP-сёрверов;
  - аварийный fallback, если tor-full окажется недоступен (manual reload).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.krab.mcp-tor</string>

    <key>Comment</key>
    <string>DEPRECATED Wave 50-B (2026-05-10) — replaced by tor-full (25 tools). Was: Krab Tor MCP SSE :8014 (Wave 44-Z).</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin/python</string>
      <string>-m</string>
      <string>src.mcp_tor_server</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/pablito/Antigravity_AGENTS/Краб</string>

    <key>StandardOutPath</key>
    <string>/tmp/krab-mcp-tor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/krab-mcp-tor.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/pablito</string>
      <key>PATH</key>
      <string>/Users/pablito/Antigravity_AGENTS/Краб/venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>MCP_TRANSPORT</key>
      <string>sse</string>
      <key>MCP_PORT</key>
      <string>8014</string>
      <key>TOR_SOCKS_PORT</key>
      <string>9050</string>
    </dict>
  </dict>
</plist>
